# standalone/pfsense-pkg-hello_world/pkg-plist
www/hello_world.php
etc/rc.d/hello_world_service.sh

%%DATADIR%%/info.xml
