# standalone/pfsense-pkg-hello_world/pkg-plist
pkg/hello_world.xml
etc/rc.d/hello_world_service.sh
www/hello_world.php

%%DATADIR%%/info.xml
