# pfBlockerNG-addons/pfsense-pkg-pfblockerng_whitelist_rulegen/pkg-plist

%%DATADIR%%/info.xml

pkg/dns_based_ip_whitelister/etc/rc.d/svc_whitelist.sh
pkg/dns_based_ip_whitelister/whitelist_monitor/__init__.py
pkg/dns_based_ip_whitelister/whitelist_monitor/__main__.py
pkg/dns_based_ip_whitelister/whitelist_monitor/app.py

pkg/dns_based_ip_whitelister/whitelist_monitor/filtersets/__init__.py
pkg/dns_based_ip_whitelister/whitelist_monitor/filtersets/common.py

pkg/dns_based_ip_whitelister/whitelist_monitor/pfsense/__init__.py
pkg/dns_based_ip_whitelister/whitelist_monitor/pfsense/api.py
pkg/dns_based_ip_whitelister/whitelist_monitor/pfsense/extract_config.php

www/dns_based_ip_whitelister/whitelist_rulegen_configure.php

pkg/dns_based_ip_whitelister.xml

pkg/dns_based_ip_whitelister/post_install_tasks.php
