# standalone/pfsense-pkg-hello_world/work/.PLIST.mktmp
/usr/local/share/licenses/hello_world-1.0/catalog.mk
/usr/local/share/licenses/hello_world-1.0/LICENSE
/usr/local/share/licenses/hello_world-1.0/APACHE20
www/hello_world.php
etc/rc.d/hello_world_service.sh

share/hello_world/info.xml
@owner
@group
@mode
